# ais_dispatcher/rootfs/etc/services.d/ais-dispatcher/run
#!/command/with-contenv bashio
set -euo pipefail

LOG_LEVEL="$(bashio::config 'log_level')"
bashio::log.info "Starting AIS Dispatcher control service"
bashio::log.info "Log level set to ${LOG_LEVEL}"

export HOME=/home/ais
export LD_LIBRARY_PATH=/home/ais/lib
export PATH="/home/ais/bin:${PATH}"
export LOG_LEVEL

bashio::log.info "Executing /home/ais/bin/aiscontrol via gosu"
if ! command -v gosu >/dev/null 2>&1; then
	bashio::log.warning "gosu command not found in PATH: ${PATH}"
fi

exec /usr/sbin/gosu ais /home/ais/bin/aiscontrol
